# vnpy_engine/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/vnpy/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:redis-wait]
command=sh -c "until redis-cli -h ${REDIS_HOST:-redis} ping > /dev/null 2>&1; do echo 'Waiting for Redis...'; sleep 2; done; echo 'Redis is ready!'; exec tail -f /dev/null"
autostart=true
autorestart=true
startretries=10
stdout_logfile=/vnpy/logs/redis-wait.log
stderr_logfile=/vnpy/logs/redis-wait.err
killasgroup=true
priority=100

[program:api-gateway]
command=python -m vnpy.api_gateway
directory=/vnpy
environment=PYTHONPATH="/vnpy",REDIS_HOST="%(env_REDIS_HOST:-redis)",REDIS_PORT="%(env_REDIS_PORT:-6379)",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=5
exitcodes=0
stdout_logfile=/vnpy/logs/api-gateway.log
stderr_logfile=/vnpy/logs/api-gateway.err
killasgroup=true
priority=200

[program:main-engine]
command=python -m vnpy.main_engine
directory=/vnpy
environment=PYTHONPATH="/vnpy",REDIS_HOST="%(env_REDIS_HOST:-redis)",REDIS_PORT="%(env_REDIS_PORT:-6379)",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=5
exitcodes=0
stdout_logfile=/vnpy/logs/main-engine.log
stderr_logfile=/vnpy/logs/main-engine.err
killasgroup=true
priority=300

[program:watchdog]
command=python -m vnpy.watchdog
directory=/vnpy
environment=PYTHONPATH="/vnpy",REDIS_HOST="%(env_REDIS_HOST:-redis)",REDIS_PORT="%(env_REDIS_PORT:-6379)",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=10
exitcodes=0
stdout_logfile=/vnpy/logs/watchdog.log
stderr_logfile=/vnpy/logs/watchdog.err
killasgroup=true
priority=400

[group:vnpy-services]
programs=api-gateway,main-engine,watchdog
priority=999

[supctl]
autostart=true
